# runner/systemd/ci-runner-manager.service
[Unit]
Description=Sanctuary CI runner reconciler
After=network-online.target libvirtd.service sanctuary-ci-firewall.service
Wants=network-online.target
Requires=libvirtd.service sanctuary-ci-firewall.service
ConditionPathExists=/etc/ci-runner/github.token

[Service]
Type=simple
User=ci-runner-manager
Group=ci-runner-manager
ExecStart=/usr/local/bin/ci-runner-manager --config /etc/ci-runner/manager.toml daemon --interval 30
LoadCredential=github_token:/etc/ci-runner/github.token
Restart=always
RestartSec=5
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/ci-runner-manager /run/lock /mnt/hdd1000-01/ci-audit
RuntimeDirectory=ci-runner-manager
RuntimeDirectoryMode=0700
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
